# deploy/tencent-cloud/systemd/zhipet-purge.timer
[Unit]
Description=Run the ZHIPET retention cleanup daily

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=10min
Persistent=true
Unit=zhipet-purge.service

[Install]
WantedBy=timers.target
